# trunk/OscarDotNet/noninteractive/.crontab
#
# $Id$
#
00 00 * * * wget -q https://sslws010.alentus.com/frompaper2web/KiAspdotnetFramework/noninteractive/milestone_sweep.aspx
