# trunk/OscarDotNet/noninteractive/.crontab
#
# $Id$
#
MAILTO=kevin.lipscomb@kvrs.org
00 00 * * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/WebEMSOF/noninteractive/milestone_sweep.aspx
15 00 * * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/OscarDotNet/noninteractive/daily_tasks.aspx
30 00 * * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/WebEMSOF_x/noninteractive/milestone_sweep.aspx
45 00 * * * wget --output-document=/dev/null https://sslws010.alentus.com/frompaper2web/OscarDotNet_x/noninteractive/daily_tasks.aspx
