# noninteractive/.crontab
#
# $Id$
#
#@daily wget -q https://sslws010.alentus.com/frompaper2web/OscarDotNet/noninteractive/milestone_sweep.aspx
@daily wget -q https://sslws010.alentus.com/frompaper2web/OscarDotNet/noninteractive/daily_tasks.aspx
